# systemd/geth.service
[Unit]
Description=Geth Ethereum Execution Client
Documentation=https://geth.ethereum.org/docs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=geth
Group=geth

ExecStart=/usr/bin/geth \
    --config /etc/geth/geth-config.toml \
    --datadir /var/lib/geth \
    --authrpc.addr localhost \
    --authrpc.port 8551 \
    --authrpc.vhosts localhost \
    --authrpc.jwtsecret /var/lib/ethereum/jwt/jwt.hex \
    --metrics \
    --metrics.addr 127.0.0.1

Restart=on-failure
RestartSec=10
TimeoutStopSec=300

LimitNOFILE=65535

PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
PrivateDevices=true

[Install]
WantedBy=multi-user.target
